# etc/profile.d/motd.sh
#!/bin/bash
clear
echo -e "\e[38;5;6m\
 .:kXWMMMWKkdolcclkMMMM:
;WMMMXx?'''        KMMM:
:MMN'              xMMM.
.WMMc             :0MMM
 dMMW;      ,     :WMMM
 .NMMWxdxkK0;     'NMMM.
  cMMMMWKx;:      'kMMM.
   :lNNl''   ,     oMMM:
          .oK;     xMMM,
      .unWMNc     .NMMd
       ':do:'     kMMk'
                .kMMx'\
\e[0m"
echo -e "\e[38;5;8m-------------------------\e[0m\n"
