# etc/systemd/system/php-fpm.service.d/privatetmp-false.conf
# allow mount for all users
[Service]
PrivateTmp=false
PrivateDevices=false

ProtectControlGroups=false
ProtectHome=false
ProtectKernelModules=false
ProtectSystem=false
