# etc/systemd/system/php-fpm.service.d/privatetmp-false.conf
[Service]
PrivateTmp=false
PrivateDevices=false

ProtectControlGroups=false
ProtectKernelModules=false
ProtectSystem=false
