# etc/systemd/system/nginx.service.d/privatetmp-false.conf
[Unit]
After=php-fpm.service

[Service]
PrivateTmp=false
