# etc/systemd/system/bluealsa.service
[Unit]
Description=Bluetooth Audio ALSA Backend
Requires=bluetooth.service
After=bluetooth.service

[Service]
ExecStart=/usr/bin/bluealsa
