# etc/systemd/system/shairport-meta.service
[Unit]
Description=Shairport Sync Metadata 
After=network.target redis.target shairport-sync.service

[Service]
ExecStart=/srv/http/bash/shairport-meta.sh
User=http
Group=http

[Install]
WantedBy=multi-user.target
