# etc/systemd/system/bluetooth.service.d/override.conf
[Service]
ExecStartPost=/srv/http/bash/system.sh bluetooth
