# etc/systemd/system/bluetooth.service.d/override.conf
[Unit]
After=nginx.service
BindsTo=bluealsa.service bluealsa-aplay.service bluez-authorize.service

[Service]
ExecStartPost=/srv/http/bash/networks.sh btset
