# etc/systemd/system/bluetooth.service.d/override.conf
[Unit]
After=nginx.service

[Service]
ExecStartPost=/usr/bin/bluetoothctl discoverable yes
ExecStartPost=/usr/bin/bluetoothctl discoverable-timeout 0
